# files/notejam.service.j2
[Unit]
Description=Unit for starting a basic Django app
After=network.service 

[Service]
Restart=on-failure
ExecStart={{ start_command }}

[Install]
WantedBy=multi-user.target
